# backend/reminder.cron
# Run daily at 08:00
0 8 * * * /app/backend/send_daily_reminders.sh >> /var/log/reminder-cron.log 2>&1
